# templates/filebeat.ini.j2
[program:filebeat]
user={{ filebeat_user }}
command={{ filebeat_install_dir }}/{{ filebeat_version }}/filebeat -c {{ filebeat_install_dir }}/{{ filebeat_version}}/filebeat.yml
autostart=true
autorestart=true
stdout_logfile={{ filebeat_log_dir }}/filebeat.out.log
stderr_logfile={{ filebeat_log_dir }}/filebeat.err.log
